# deploy/systemd/trading-assistant-backup.timer
[Unit]
Description=Daily Trading Assistant backup

[Timer]
OnCalendar=*-*-* 03:15:00 UTC
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target

